# cron/backup.cron
0 2 * * * /home/vagrant/Linux-Backup-restore-system/backup.sh
 

